# email/offlineimap.service
[Unit]
Description=Offline IMAP syncing

[Service]
ExecStart=/usr/bin/offlineimap

[Install]
WantedBy=default.target
