# email/offlineimap-personal.service
[Unit]
Description=Offline IMAP syncing

[Service]
ExecStart=/usr/bin/offlineimap -s -c ~/.offlineimaprc.personal

[Install]
WantedBy=default.target
